# controller/configuration/web.service
[Service]
Type=simple
ExecStart=/usr/local/bin/node /var/www/organization/web/index.js
Restart=on-failure

[Install]
WantedBy=multi-user.target
